# supabase/functions/_cron/wallet-push.cron
# Process Apple Wallet push notifications every minute
# This triggers the wallet-push edge function to send pending notifications
*/1 * * * * wallet-push
